# conf/supervisord-ssh.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/bin/sshd -D
autorestart=true

