# conf/supervisord-ipython.ini
[program:ipython]
directory = /home/ipy/ipynotebooks/
user = ipy
environment = IPYTHONDIR=/home/ipy/.ipython
command=ipython2 notebook --profile=nbserver
stopsignal = INT
killasgroup = true
autostart = true
autorestart = false
redirect_stderr = true
